# supervisor.ini
[program:phalcon]
command                 = php run queue
directory               = /path/to/app
process_name            = %(program_name)s
numprocs                = 1
autostart               = true
autorestart             = true
stdout_logfile          = /path/to/app/storage/log/stdout.log
stdout_logfile_maxbytes = 10MB
stderr_logfile          = /path/to/app/storage/log/stderr.log
stderr_logfile_maxbytes = 10MB